# examples/systemd/linux-healthy-agent.service
[Unit]
Description=Linux Healthy Agent
Documentation=https://github.com/lauzhihao/linux-healthy-agent-with-feishu-hooks

[Service]
Type=oneshot
EnvironmentFile=-/etc/linux-healthy-agent.env
ExecStart=/usr/local/bin/linux-healthy-agent --alert-state-file /run/linux-healthy-agent-alert.json
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/run
